# crontab
* * * * * cd /var/www/woltlab && /usr/local/bin/php cron.php > /dev/null
# An empty line is required at the end of this file for a valid cron file.